# flux/data/Library/LaunchDaemons/org.herf.flux.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>org.herf.flux</string>
	<key>UserName</key>
	<string>root</string>
	<key>Program</key>
	<string>/usr/libexec/flux</string>
	<key>JetsamProperties</key>
	<dict>
		<key>JetsamPriority</key>
		<integer>-49</integer>
		<key>JetsamMemoryLimit</key>
		<integer>10240</integer>
	</dict>
	<key>KeepAlive</key>
	<true/>
	<key>EnablePressuredExit</key>
	<false/>
	<key>RunAtLoad</key>
	<true/>
</dict>
</plist>
